# gunicorn.socket
[Unit]
Description=gunicorn socket for Django Proyecto-Lenguaje-Gestos

[Socket]
# Socket Unix para comunicación con Nginx
ListenStream=/run/gunicorn.sock

# Permisos del socket
# root puede crear el socket, www-data puede acceder
SocketUser=root
SocketGroup=www-data
SocketMode=0660

# Configuración del socket
Backlog=512
NoDelay=true

[Install]
WantedBy=sockets.target


